# chefknowledge/site-cookbooks/k8s/templates/default/k8s-node.service.erb
[Unit]
Description=Kubernetes Minion

[Service]
Restart=on-failure 
RestartSec=5s
ExecStart=/usr/local/bin/hyperkube kubelet --api_servers=<%= @API_SERVER %> --v=2 --address=0.0.0.0 --enable-server


[Install]
WantedBy=multi-user.target 

